# files/etc/uci-defaults/99_gameboost
#!/bin/sh

# Initialize default UCI config for gameboost

CONFIG_FILE="/etc/config/gameboost"

if [ ! -f "$CONFIG_FILE" ]; then
    touch "$CONFIG_FILE"
fi

uci -q batch <<-EOF >/dev/null
    set gameboost.global=gameboost
    set gameboost.global.enabled='0'
    set gameboost.global.rules_source='https://raw.githubusercontent.com/GameBoost-rules/main/rules.txt'
    set gameboost.global.auto_update='1'
    set gameboost.global.update_interval='86400'
    set gameboost.global.log_level='info'
    set gameboost.global.proxy_enabled='0'      # 预留 iptables 代理开关
    set gameboost.global.proxy_port='1080'      # 代理监听端口
EOF

# 预设平台选项
for platform in steam epic origin uplay battle net xbox psn nswitch; do
    uci set gameboost.${platform}=platform
    uci set gameboost.${platform}.enabled='1'
done

uci commit gameboost

exit 0
